# sinopia.ini
[program:sinopia]
command=sinopia
directory=/usr/bin/
autostart=true
autorestart=true
numprocs=1

