# sinopia.ini
[program:sinopia]
command=node /usr/local/node/bin/sinopia
directory=/usr/bin/local/node/bin
autostart=true
autorestart=true
numprocs=1

