# verdaccio.ini
[program:verdaccio]
command=node /usr/local/node/bin/verdaccio -c /root/.config/verdaccio/config.yaml
#directory=/usr/bin/local/node/bin
autostart=true
autorestart=true
numprocs=1

